# ops/crontab/crontab.template
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin

# ai-basecamp: git pull + npm run build + goaccess
0,30 * * * * /home/dukai/WorkSpace/ai-basecamp/ops/cron-ai-basecamp.sh
